# system/firstboot/starfleet-firstboot.service

[Unit]
Description=Starfleet first-boot finisher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/starfleet-firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
